# deploy/systemd/dingpan-generate.service
[Unit]
Description=DingPan Generate Shared Analysis
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/dingpan
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/dingpan/.venv/bin/python /var/www/dingpan/generate.py --force --today-if-trading-day
TimeoutStartSec=10min

[Install]
WantedBy=multi-user.target
